# ex04/81-mndlovu-usb.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+=/sbin/insmod /home/mndlovu/kernel/little_penguin/ex04/main.ko
